# deploy/lifecoach.service
[Unit]
Description=Life Coach Recording Assistant
After=network.target

[Service]
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/LifeCoach
Environment="PATH=/home/sunrise/LifeCoach/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/bin/bash /home/sunrise/LifeCoach/deploy/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
